# cmd/lima-guestagent/lima-guestagent.TEMPLATE.service
[Unit]
Description=lima-guestagent

[Service]
ExecStart={{.Binary}} daemon
Type=simple
Restart=on-failure

[Install]
WantedBy=default.target
